# examples/displayboard.service
[Unit]
Description=Displayboard Python Application
# Don't require network - app can run offline after bootstrap
After=local-fs.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/displayboard/current

# Run using the Poetry venv
# First try .venv in project, then use full path to poetry to find venv
ExecStart=/bin/bash -c '\
  if [ -d .venv ]; then \
    exec .venv/bin/python -m displayboard.main -d; \
  elif [ -x $HOME/.local/bin/poetry ]; then \
    VENV_PATH=$($HOME/.local/bin/poetry env info --path 2>/dev/null); \
    exec $VENV_PATH/bin/python -m displayboard.main -d; \
  else \
    echo "Error: Cannot find Poetry venv" >&2; \
    exit 1; \
  fi'

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=displayboard

# Security and resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/displayboard /var/log/displayboard

[Install]
WantedBy=multi-user.target
