# examples/displayboard.service
[Unit]
Description=Displayboard Python Application
# Don't require network - app can run offline after bootstrap
After=local-fs.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/displayboard/current

# Run using the Poetry venv
# Poetry creates venv in project directory as .venv (if configured) or in cache dir
# This uses the .venv if present, otherwise tries to find it via poetry
ExecStart=/bin/bash -c '\
  if [ -d .venv ]; then \
    exec .venv/bin/python -m displayboard.main -d; \
  else \
    VENV_PATH=$(poetry env info --path 2>/dev/null || echo ".venv"); \
    exec $VENV_PATH/bin/python -m displayboard.main -d; \
  fi'

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=displayboard

# Security and resource limits
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/displayboard /var/log/displayboard

[Install]
WantedBy=multi-user.target
