# examples/gh-deployer.service
[Unit]
Description=GitHub Release Deployer for Displayboard
After=network-online.target
Wants=network-online.target
# Wait for network to actually be available before starting
Requires=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/displayboard/gh-deployer

# Run the deployer with the config file
ExecStart=/usr/local/bin/gh-deployer --config /opt/displayboard/gh-deployer/config.yaml

# Restart policy - keep trying if it fails or network drops
Restart=always
RestartSec=60
# Start immediately on system boot once network is available
StartLimitIntervalSec=0

# Environment variables
Environment="GITHUB_TOKEN_FILE=/opt/displayboard/config/github-token"

# Logging (will show network retry attempts)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gh-deployer

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/displayboard /var/log/displayboard

[Install]
WantedBy=multi-user.target
