# examples/gh-deployer.service
[Unit]
Description=GitHub Release Deployer for Displayboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/displayboard/gh-deployer

# Run the deployer with the config file
ExecStart=/usr/local/bin/gh-deployer --config /opt/displayboard/gh-deployer/config.yaml

# Restart policy - keep trying if it fails
Restart=always
RestartSec=30

# Environment variables
Environment="GITHUB_TOKEN_FILE=/opt/displayboard/config/github-token"

# Read GitHub token from file if it exists
ExecStartPre=/bin/bash -c 'if [ -f /opt/displayboard/config/github-token ]; then export GITHUB_TOKEN=$(cat /opt/displayboard/config/github-token); fi'

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gh-deployer

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/displayboard /var/log/displayboard

[Install]
WantedBy=multi-user.target
